# deluge/1.3/fs/etc/supervisor.d/supervisor.ini
[supervisord]
nodaemon=true

[program:deluged]
user=torrent
environment=TERM="xterm"
command=/usr/bin/deluged -d -c /config -L info -l /config/deluged.log

[program:deluge-web]
user=torrent
environment=PYTHON_EGG_CACHE="/home/torrent/.python-eggs",TERM="xterm"
command=/usr/bin/deluge-web -c /config -L info -l /config/deluge-web.log
