# subsonic/old-5.3/supervisord.conf
[supervisord]
nodaemon=true

[program:subsonic]
user=subsonic
environment=TERM="xterm"
command=./opt/subsonic/subsonic.sh
