# emby/rootfs/etc/supervisor.d/emby.ini
[supervisord]
nodaemon=true

[program:emby]
user=emby
directory=/embyServer
priority=2
redirect_stderr=true
environment=HOME="/embyData",PWD="/embyServer",LOGNAME="emby",USER="emby",TERM="xterm"
startsecs=5
autorestart=true
command=mono MediaBrowser.Server.Mono.exe -ffmpeg /usr/bin/ffmpeg -ffprobe /usr/bin/ffprobe
